# roles/macos-general-config/templates/mlx-bridge.plist.j2
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.thadigus.mlx-bridge</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{ dotfiles_dir }}/scripts/mlx_bridge.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>{{ user_home }}/Library/Logs/mlx-bridge.log</string>
    <key>StandardErrorPath</key>
    <string>{{ user_home }}/Library/Logs/mlx-bridge.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>MLX_PORT</key>
        <string>{{ mlx_port }}</string>
        <key>MLX_WG_PREFIX</key>
        <string>{{ mlx_wg_prefix }}</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/bin:/bin:/sbin</string>
    </dict>
</dict>
</plist>
